# docker/vault_supervisord.ini
[program:consul]
priority=1
user=consul
command=/usr/local/consul -config-file=/etc/consul.d/config.json -pid-file=/etc/consul.d/consul.pid
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout

[program:vault-0.9.6-dev]
priority=2
user=vault-0.9.6-dev
command=/usr/local/bin/vault.0.9.6 -dev -dev-listen-address=127.0.0.1:8200 -dev-root-token-id=${VAULT-0.9.6-DEV} -pid-file=/etc/vault.d-0.9.6-dev/vault.pid
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout

[program:vault-1.5.4-dev]
priority=2
user=vault-1.5.4-dev
command=/usr/local/bin/vault-1.5.4 -dev -dev-listen-address=127.0.0.1:8400 -dev-root-token-id=${VAULT-1.5.4-DEV} -pid-file=/etc/vault.d-1.5.4-dev/vault.pid
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout

[program:vault-0.9.6-prod]
priority=3
user=vault-0.9.6-prod
command=/usr/local/bin/vault.0.9.6 -config=/etc/vault.d-0.9.6-prod/config.hcl
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout

[program:vault-1.5.4-prod]
priority=3
user=vault-1.5.4-prod
command=/usr/local/bin/vault-1.5.4 -config=/etc/vault.d-1.5.4-prod/config.hcl
stderr_logfile = /dev/stderr
stdout_logfile = /dev/stdout
